# supervisord.conf
[supervisord]
nodaemon=true

[program:mailcatcher]
command=/usr/local/bin/mailcatcher --ip 0.0.0.0 --smtp-port %(ENV_SMTP_PORT)s --http-port %(ENV_HTTP_PORT)s --foreground
autorestart=true
